# systemd/tbaps-backend.service
[Unit]
Description=TBAPS Backend API (FastAPI + Uvicorn)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tbaps/backend
ExecStart=/opt/tbaps/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
EnvironmentFile=/opt/tbaps/.env
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
